# pi-5/etc/profile.d/09-splashscreen-wait.sh
while pgrep mpv >/dev/null; do sleep 1; done #Splashscreen-wait
